# deployment-files/unicorn.conf
[program:unicorn]
command = /var/www/crucible_smart_app/deployment-files/start_unicorn.sh    ; Command to start app
user = ubuntu                                                          ; User to run as
stdout_logfile = /var/www/crucible_smart_app/tmp/unicorn_supervisor.log   ; Where to write log messages
redirect_stderr = true                                                ; Save stderr in the same log
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8 ; Set UTF-8 as default encoding
